# configs/user_custom/egpu/nvidia-egpu.service
[Unit]
Description=NVIDIA eGPU Setup
After=bolt.service
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nvidia-egpu-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
